# deploy/wedding-app.service.example
[Unit]
Description=Himmet & Cennet düğün yükleme (FastAPI)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/wedding-app
Environment="PATH=/var/www/wedding-app/venv/bin"
ExecStart=/var/www/wedding-app/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8000 --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
